# system-config/keyswap.service
[Unit]
Description=keyswap for Terminal and GUI
PartOf=graphical-session.target

[Service]
Type=simple
Restart=always
#Environment=DISPLAY=:0.0
RestartSec=1
WorkingDirectory=/home/{username}/.config/kinto
ExecStart=/bin/bash -c "/home/{username}/.config/kinto/xactive.sh carrots"
ExecStop=/bin/bash /home/{username}/.config/kinto/cleanup.sh

[Install]
WantedBy=default.target
